# usr/lib/systemd/system/spice-vdagentd.service.d/30_vm-config-dist.conf
## Copyright (C) 2020 - 2025 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

## SPICE is a QEMU/KVM technology; VirtualBox does not implement it. Do not run
## the SPICE guest agent daemon under VirtualBox (belt-and-suspenders alongside
## the matching drop-in on spice-vdagentd.socket).

[Unit]
ConditionVirtualization=!oracle
