# usr/lib/systemd/system/spice-vdagentd.service.d/30_vm-config-dist.conf
## Copyright (C) 2020 - 2025 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

## VirtualBox does not support SPICE: https://www.virtualbox.org/ticket/13042
[Unit]
ConditionVirtualization=kvm
## Or, if we support QEMU with hypervisors other than KVM in the future:
# ConditionVirtualization=|kvm
# ConditionVirtualization=|qemu
