# usr/lib/systemd/user/spice-vdagent.service.d/30_vm-config-dist.conf
## Copyright (C) 2020 - 2025 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

## SPICE is a QEMU/KVM technology; VirtualBox does not implement it. Do not run
## the SPICE guest session agent under VirtualBox, where it has no daemon to
## connect to.

[Unit]
ConditionVirtualization=!oracle
