# usr/lib/systemd/system/vboxadd-service.service.d/30_vm-config-dist.conf
## Copyright (C) 2012 - 2025 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

## VirtualBox Guest Additions installs this unit generated from an LSB init
## script, without a virtualization condition. Only run it under VirtualBox.

[Unit]
ConditionVirtualization=oracle
