# etc/xdg/autostart/wlr-resize-watcher.desktop
[Desktop Entry]
Type=Application
Name=Virtual display resize helper
Exec=/usr/bin/wlr-resize-watcher
StartupNotify=false
NoDisplay=true
