# usr/lib/systemd/system/spice-vdagentd.socket.d/30_vm-config-dist.conf
## Copyright (C) 2020 - 2025 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

## SPICE is a QEMU/KVM technology; VirtualBox does not implement it. Under
## VirtualBox this socket still listens and activates spice-vdagentd, which
## then has no SPICE host to talk to (confirmed running on an oracle guest).
## Skip the activation socket under VirtualBox so the daemon is never started.

[Unit]
ConditionVirtualization=!oracle
